# com.anwarselo.dart-taskmaster.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anwarselo.dart-taskmaster</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/anwarselo/Desktop/Development/test/start-integration.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/anwarselo/Desktop/Development/test/integrations/dart-taskmaster/logs/launchd-error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/anwarselo/Desktop/Development/test/integrations/dart-taskmaster/logs/launchd-out.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
